# etc/systemd/backup.service
[Unit]
Description=Pi backup

[Service]
ExecStart=/mnt/internal1/Scripts/rsync/backup.sh

[Install]
WantedBy=default.target
