# etc/systemd/backup.path
[Unit]
Description=Check for changes and perform backup

[Path]
PathChanged=/mnt/pi/pi-os
